# lib/profile.d/psql.sh
export PSQL_HISTORY=~/.cache/pg/psql_history
