# lib/profile.d/ruby.sh
export BUNDLE_USER_CONFIG=~/.config/bundle
export BUNDLE_USER_CACHE=~/.cache/bundle
export BUNDLE_USER_PLUGIN=~/.local/share/bundle
