# dot.config/systemd/user/fontsrv.service
[Unit]
Description=fontsrv
After=default.target

[Service]
Type=simple
Environment=PLAN9=/usr/local/plan9
ExecStart=/usr/local/plan9/bin/fontsrv

[Install]
WantedBy=default.target
